# nomad/.config/autostart/dsbautostart-LSE1Dz.desktop
[Desktop Entry]
Type=Application
Exec=dsbbatmon
Hidden=false
Terminal=false
