# nomad/.config/autostart/dsbautostart-j6wQfu.desktop
[Desktop Entry]
Type=Application
Exec=tint2
Hidden=false
Terminal=false
