# nomad/.config/autostart/dsbautostart-00uGIR.desktop
[Desktop Entry]
Type=Application
Exec=~/.config/DSB/dsbxinput.sh
Hidden=false
Terminal=false
