# nomad/.config/autostart/dsbautostart-OByvDf.desktop
[Desktop Entry]
Type=Application
Exec=dsbmc -i
Hidden=false
Terminal=false
