# nomad/.config/autostart/dsbautostart-yEsPCp.desktop
[Desktop Entry]
Type=Application
Exec=brightness_listener
Hidden=false
Terminal=false
