# nomad/.config/autostart/dsbautostart-jub5kH.desktop
[Desktop Entry]
Type=Application
Exec=chkmultihead
Hidden=false
Terminal=false
