# nomad/.config/autostart/dsbautostart-j3p2Mr.desktop
[Desktop Entry]
Type=Application
Exec=plank
Hidden=false
Terminal=false
