# nomad/.config/autostart/dsbautostart-6YexdA.desktop
[Desktop Entry]
Type=Application
Exec=xbindkeys
Hidden=false
Terminal=false
