# nomad/.config/autostart/dsbautostart-rNGiiz.desktop
[Desktop Entry]
Type=Application
Exec=sudo -E networkmgr
Hidden=false
Terminal=false
