# nomad/.config/autostart/dsbautostart-KFClXF.desktop
[Desktop Entry]
Type=Application
Exec=~/.config/DSB/dsbds.sh
Hidden=false
Terminal=false
