# nomad/.config/autostart/dsbautostart-KFClXF.desktop
[Desktop Entry]
Type=Application
Exec=sh -c '~/.config/DSB/dsbds.sh; ~/.fehbg'
Hidden=false
Terminal=false
