# nomad/.config/autostart/dsbautostart-bo2rzY.desktop
[Desktop Entry]
Type=Application
Exec=~/.screenlayout/default.sh
Hidden=false
Terminal=false
