# nomad/.config/autostart/dsbautostart-yWjNU2.desktop
[Desktop Entry]
Type=Application
Exec=dsbmixer -i
Hidden=false
Terminal=false
