# nomad/.config/autostart/dsbautostart-xF3CBL.desktop
[Desktop Entry]
Type=Application
Exec=picom -b
Hidden=false
Terminal=false
