# nomad/.config/autostart/dsbautostart-d7HlBM.desktop
[Desktop Entry]
Type=Application
Exec=redshift -l$(getlalo)
Hidden=false
Terminal=false
