# nomad/.config/autostart/dsbautostart-xF3CBL.desktop
[Desktop Entry]
Type=Application
Exec=compton
Hidden=false
Terminal=false
