# nomad/.config/autostart/dsbautostart-d3wjFP.desktop
[Desktop Entry]
Type=Application
Exec=dunst
Hidden=false
Terminal=false
