# nomad/.config/autostart/hplip-systray.desktop
[Desktop Entry]
Type=Application
Name=HP System Tray Service
Comment=HP System Tray Service
Exec=hp-systray -x
Hidden=true
Terminal=false
