# nomad/.config/autostart/ibus.desktop
[Desktop Entry]
Type=Application
Name=IBus
Comment=Start IBus Input Method Framework
Exec=ibus-daemon -drxR
Hidden=false
Terminal=false
