# config/etc/rc.d/mount_uzip
#!/bin/sh
#
# PROVIDE: mount_uzip
# BEFORE: ldconfig NETWORKING
# REQUIRE: FILESYSTEMS

. /etc/rc.subr 

echo "Mounting uzip. Please wait ..."
mdmfs -P -F /uzip/usr.local.uzip -oro md.uzip /usr/local

